# ansible/roles/minecraft_server/templates/minecraft.service.j2
[Unit]
Description=Minecraft Server
After=network.target

[Service]
User=minecraft
Group=minecraft
WorkingDirectory=/opt/minecraft
ExecStart=/usr/bin/screen -dmS minecraft /usr/bin/java -Xms2G -Xmx4G -jar server.jar nogui
ExecStop=/usr/bin/screen -S minecraft -X stuff "stop\n"
TimeoutStopSec=60
Restart=on-failure
RestartSec=5
Type=forking
StandardOutput=journal
StandardError=journal
KillMode=mixed

[Install]
WantedBy=multi-user.target